# linux/systemd/drivertranslator.service
[Unit]
Description=RTI DriverTranslator (WyreStorm NHD-CTL emulator -> AMX)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/drivertranslator
ExecStart=/opt/drivertranslator/.venv/bin/python -m drivertranslator --config /opt/drivertranslator/config.json --listen 0.0.0.0 --port 2323
Restart=always
RestartSec=2

# Hardening (safe defaults; adjust if needed)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target

